# packaging/server/conf/init_web2py_scheduler.conf
description \"web2py vidjil task scheduler\"
start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown
respawn
respawn limit 8 60 # Give up if restart occurs 8 times in 60 seconds.
exec  sudo -u www-data python /usr/share/vidjil/server/web2py/web2py.py -K vidjil,vidjil,vidjil
