# docker/reporter/crontab
0 5 * * * root /usr/bin/python /opt/reporter/reporter.py >> /var/log/cron.log 2>&1

